# school-portal/docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/www/html/storage/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[program:laravel-worker-default]
command=php /var/www/html/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --timeout=60
directory=/var/www/html
autostart=true
autorestart=true
startretries=3
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/www/html/storage/logs/worker-default.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/www/html/storage/logs/worker-default-error.log
stderr_logfile_maxbytes=5MB
stopwaitsecs=3600

[program:laravel-worker-notifications]
command=php /var/www/html/artisan queue:work redis --queue=notifications --sleep=3 --tries=3 --max-time=3600 --timeout=30
directory=/var/www/html
autostart=true
autorestart=true
startretries=3
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/www/html/storage/logs/worker-notifications.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/www/html/storage/logs/worker-notifications-error.log
stderr_logfile_maxbytes=5MB

[program:laravel-worker-media]
command=php /var/www/html/artisan queue:work redis --queue=media --sleep=3 --tries=3 --max-time=3600 --timeout=120
directory=/var/www/html
autostart=true
autorestart=true
startretries=3
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/www/html/storage/logs/worker-media.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/www/html/storage/logs/worker-media-error.log
stderr_logfile_maxbytes=5MB

[program:laravel-scheduler]
command=bash -c "while true; do php /var/www/html/artisan schedule:run >> /var/www/html/storage/logs/scheduler.log 2>&1; sleep 60; done"
directory=/var/www/html
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/www/html/storage/logs/scheduler.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/www/html/storage/logs/scheduler-error.log
stderr_logfile_maxbytes=5MB
